# server-name-v1.0/DEBIAN/prerm
#!/bin/bash
systemctl stop prometheus
systemctl disable prometheus
ufw delete allow in on hoofhouse from 192.168.128.0/17 to any port 9090 comment prometheus
